# container/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/usr/share/logstash/server

[program:logstash]
command=/usr/share/logstash/bin/logstash
